# debian-openclaw/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        # Create openclaw user/group if they don't exist
        if ! getent group openclaw >/dev/null 2>&1; then
            groupadd --system openclaw
        fi
        if ! getent passwd openclaw >/dev/null 2>&1; then
            useradd --system --gid openclaw --home-dir /var/lib/openclaw \
                    --shell /usr/sbin/nologin openclaw
        fi

        # Ensure state directories exist with correct permissions
        mkdir -p /var/lib/openclaw
        chown openclaw:openclaw /var/lib/openclaw
        chmod 750 /var/lib/openclaw

        mkdir -p /var/log/openclaw
        chown openclaw:openclaw /var/log/openclaw
        chmod 750 /var/log/openclaw

        # Copy example config if no config exists
        if [ ! -f /etc/openclaw/openclaw.json ]; then
            if [ -f /etc/openclaw/openclaw.json.example ]; then
                cp /etc/openclaw/openclaw.json.example /etc/openclaw/openclaw.json
                chmod 640 /etc/openclaw/openclaw.json
                chown root:openclaw /etc/openclaw/openclaw.json
            fi
        fi

        # Load kernel module if available
        if [ -f /lib/modules/$(uname -r)/extra/openclaw.ko ]; then
            modprobe openclaw 2>/dev/null || true
        fi

        # Enable and start the gateway service
        if command -v systemctl >/dev/null 2>&1; then
            systemctl daemon-reload
            systemctl enable openclawgw.service 2>/dev/null || true
            systemctl start openclawgw.service 2>/dev/null || true
        fi

        echo "OpenClaw installed successfully."
        echo "  CLI:     openclaw agent -m 'Hello from C'"
        echo "  Config:  /etc/openclaw/openclaw.json"
        echo "  Gateway: systemctl status openclawgw"
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument '$1'" >&2
        exit 1
        ;;
esac

exit 0
