# debian-openclaw/prerm
#!/bin/sh
set -e

case "$1" in
    remove|upgrade|deconfigure)
        # Stop the gateway service
        if command -v systemctl >/dev/null 2>&1; then
            systemctl stop openclawgw.service 2>/dev/null || true
            systemctl disable openclawgw.service 2>/dev/null || true
        fi

        # Unload kernel module
        if lsmod | grep -q "^openclaw "; then
            rmmod openclaw 2>/dev/null || true
        fi
        ;;

    failed-upgrade)
        ;;

    *)
        echo "prerm called with unknown argument '$1'" >&2
        exit 1
        ;;
esac

exit 0
